# services/crio_node/crio_node.service
[Unit]
Description=NISystem cRIO Node Service
Documentation=https://github.com/your-repo/nisystem
After=network-online.target mosquitto.service
Wants=network-online.target

[Service]
Type=simple
User=admin
Group=ni

# Environment file for MQTT broker configuration
EnvironmentFile=-/home/admin/nisystem/crio_node.env

# Working directory
WorkingDirectory=/home/admin/nisystem

# Start the service
ExecStart=/usr/bin/python3 /home/admin/nisystem/crio_node.py --broker ${MQTT_BROKER} --port ${MQTT_PORT} --node-id ${NODE_ID}

# Restart on failure
Restart=always
RestartSec=5

# Give it time to cleanup on stop
TimeoutStopSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=crio-node

[Install]
WantedBy=multi-user.target
